# pkg/50_zellij/local.zellij-theme-switcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.zellij-theme-switcher</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <string>set-zellij-colorscheme</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>  <!-- run every 60 seconds -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/zellij-theme-switcher.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/zellij-theme-switcher.err</string>
</dict>
</plist>
