# roles/common/files/dotfiles/.config/systemd/user/backup-sync.timer
[Unit]
Description=Daily backup of Sync folder

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
