# roles/common/files/dotfiles/.config/systemd/user/backup-sync.service
[Unit]
Description=Backup Sync folder

[Service]
Type=oneshot
ExecStart=/home/jlrickert/Sync/dotfiles/bin/backup-sync
